# DrinkWater/drink-water.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin/:/root/bin
MAILTO=root
HOME=/

01 09-20 * * 1-5 root /bin/bash /DRINK_WATER_DIR/DrinkWater.sh
