# Dict/dict.cron
59 01 * * * root sh /DICT_PLUGIN_DIR/DictSendEmail.sh
